# os/linux/autostart/xmodmap.desktop
[Desktop Entry]
Name=pgupdn
Exec=xmodmap ~/.Xmodmap
Type=Application